# web/nginx-s6-overlay/rootfs/etc/services.d/crystal-api-3/run
#!/bin/bash
echo "Starting Crystal API instance 3 on port 3003..."
export PORT=3003
export INSTANCE_ID=3
exec /usr/bin/my-api
