# web/nginx-s6-overlay/rootfs/etc/services.d/crystal-api-2/run
#!/bin/bash
echo "Starting Crystal API instance 2 on port 3002..."
export PORT=3002
export INSTANCE_ID=2
exec /usr/bin/my-api
