# web/nginx-s6-overlay/rootfs/etc/services.d/crystal-api-1/run
#!/bin/bash
echo "Starting Crystal API instance 1 on port 3001..."
export PORT=3001
export INSTANCE_ID=1
exec /usr/bin/my-api
